# src/installboot/gnome-ostree-install-bootloader.target
[Unit]
Description=Bootloader installer
Requires=gnome-ostree-install-bootloader.service
AllowIsolate=yes
