# src/gnome-ostree-initial-setup-by-default.service.in
[Unit]
Description=Enable Initial Setup
ConditionPathExists=!/var/lib/gdm/completed-initial-setup

[Service]
ExecStart=@pkglibdir@/initial-setup-by-default
Type=oneshot
Before=gdm.service
